# srcs/mysql/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug
logfile=supervisord.log
pidfile=supervisord.pid

[program:telegraf]
command= telegraf
autostart=true
autorestart=true

[program:mysql]
command= mysqld_safe --datadir="/var/lib/mysql"
autostart=true
autorestart=true